# inlab_plasmic/.build/supervisord.conf
# README

# https://www.digitalocean.com/community/tutorials/how-to-install-and-manage-supervisor-on-ubuntu-and-debian-vps
# http://supervisord.org/configuration.html#program-x-section-settings
# https://docs.docker.com/config/containers/multi-service_container/#use-a-process-manager

[unix_http_server]
file=/tmp/supervisor.sock
chmod=0755

[supervisord]
nodaemon=true
pidfile=/tmp/supervisord.pid
# logfile=/dev/null
# logfile_maxbytes=0

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:npm]
command=npm run dev
directory=/app/
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
stdout_events_enabled=true
stderr_events_enabled=true
